# services/controller/run
#!/bin/bash

set -e
set -u

cd "$(dirname "${BASH_SOURCE[0]}")"

export PYTHONPATH=./src

redis-server &

exec gunicorn \
  --bind="0.0.0.0:${PORT:-8080}" \
  --preload \
  --workers=1 \
  --timeout=120 \
  --pythonpath=${PYTHONPATH} \
  plz.controller.main:app \
  -- $@
