# services/controller/run
#!/bin/bash

export PYTHONPATH=./src

exec python ./src/plz/controller/main.py "$@"
