# services/controller/run
#!/bin/bash

set -e
set -u

cd "$(dirname "${BASH_SOURCE[0]}")"

if [[ -z "${CONFIGURATION:-}" ]]; then
  export CONFIGURATION="$(cat "$1")"
fi
PORT="${PORT:-$(python -c 'import os; import pyhocon; print(pyhocon.ConfigFactory.parse_string(os.environ.get("CONFIGURATION", "")).get("port", ""))')}"
PORT="${PORT:-8080}"

export PYTHONPATH='./src'

exec gunicorn \
  --bind="0.0.0.0:${PORT}" \
  --preload \
  --workers=1 \
  --timeout=120 \
  --pythonpath="${PYTHONPATH}" \
  plz.controller.main:app \
  -- $@
