# services/controller/run
#!/bin/bash

set -e
set -u

cd "$(dirname "${BASH_SOURCE[0]}")"

export PYTHONPATH=./src
exec python ./src/plz/controller/main.py "$@"
